# next-app-template/com.uspto.auto-download.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.uspto.auto-download</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/etiennetabchoury/Documents/JaypeeProjects/Shefle-project/next-app-template &amp;&amp; npx tsx scripts/uspto/download_latest_trtdxfap.ts >> logs/uspto-auto.log 2>&1</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/etiennetabchoury/Documents/JaypeeProjects/Shefle-project/next-app-template/logs/uspto-auto.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/etiennetabchoury/Documents/JaypeeProjects/Shefle-project/next-app-template/logs/uspto-error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
